# data/rules/cyber_security/mitre_attack/json/Operating System Configuration.txt.json
[
    "Enterprise | T1548 | Abuse Elevation Control Mechanism | Applications with known vulnerabilities or known shell escapes should not have the setuid or setgid bits set to reduce potential damage if an application is compromised. Additionally, the number of programs with setuid or setgid bits set should be minimized across a system. Ensuring that the sudo tty_tickets setting is enabled will prevent this leakage across tty sessions.",
    " |  | .001 | Setuid and Setgid | Applications with known vulnerabilities or known shell escapes should not have the setuid or setgid bits set to reduce potential damage if an application is compromised. Additionally, the number of programs with setuid or setgid bits set should be minimized across a system.",
    " |  | .003 | Sudo and Sudo Caching | Ensuring that thetty_ticketssetting is enabled will prevent this leakage across tty sessions.",
    "Enterprise | T1087 | Account Discovery | Prevent administrator accounts from being enumerated when an application is elevating through UAC since it can lead to the disclosure of account names. The Registry key is locatedHKLM\\ SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\CredUI\\EnumerateAdministrators. It can be disabled through GPO: Computer Configuration > [Policies] > Administrative Templates > Windows Components > Credential User Interface: E numerate administrator accounts on elevation.[1]",
    " |  | .001 | Local Account | Prevent administrator accounts from being enumerated when an application is elevating through UAC since it can lead to the disclosure of account names. The Registry key is located atHKLM\\ SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\CredUI\\EnumerateAdministrators. It can be disabled through GPO: Computer Configuration > [Policies] > Administrative Templates > Windows Components > Credential User Interface: Enumerate administrator accounts on elevation.[1]",
    " |  | .002 | Domain Account | Prevent administrator accounts from being enumerated when an application is elevating through UAC since it can lead to the disclosure of account names. The Registry key is located atHKLM\\ SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\CredUI\\EnumerateAdministrators. It can be disabled through GPO: Computer Configuration > [Policies] > Administrative Templates > Windows Components > Credential User Interface: Enumerate administrator accounts on elevation.[1]",
    "Enterprise | T1098 | Account Manipulation | Protect domain controllers by ensuring proper security configuration for critical servers to limit access by potentially unnecessary protocols and services, such as SMB file sharing.",
    "Enterprise | T1197 | BITS Jobs | Consider reducing the default BITS job lifetime in Group Policy or by editing theJobInactivityTimeoutandMaxDownloadTimeRegistry values inHKEY_LOCAL_MACHINE\\Software\\Policies\\Microsoft\\Windows\\BITS.[2]",
    "Enterprise | T1092 | Communication Through Removable Media | Disallow or restrict removable media at an organizational policy level if they are not required for business operations.[3]",
    "Enterprise | T1136 | Create Account | Protect domain controllers by ensuring proper security configuration for critical servers.",
    " |  | .002 | Domain Account | Protect domain controllers by ensuring proper security configuration for critical servers.",
    "Enterprise | T1543 | Create or Modify System Process | Ensure that Driver Signature Enforcement is enabled to restrict unsigned drivers from being installed.",
    " |  | .003 | Windows Service | Ensure that Driver Signature Enforcement is enabled to restrict unsigned drivers from being installed.",
    "Enterprise | T1546 | .008 | Event Triggered Execution:Accessibility Features | To use this technique remotely, an adversary must use it in conjunction with RDP. Ensure that Network Level Authentication is enabled to force the remote desktop session to authenticate before the session is created and the login screen displayed. It is enabled by default on Windows Vista and later.[4]",
    "Enterprise | T1011 | Exfiltration Over Other Network Medium | Prevent the creation of new network adapters where possible.[5][6]",
    " |  | .001 | Exfiltration Over Bluetooth | Prevent the creation of new network adapters where possible.",
    "Enterprise | T1564 | .002 | Hide Artifacts:Hidden Users | If the computer is domain joined, then group policy can help restrict the ability to create or hide users. Similarly, preventing the modification of the/Library/Preferences/com.apple.loginwindowHide500Usersvalue will force all users to be visible.",
    "Enterprise | T1574 | .006 | Hijack Execution Flow:Dynamic Linker Hijacking | When System Integrity Protection (SIP) is enabled in macOS, the aforementioned environment variables are ignored when executing protected binaries. Third-party applications can also leverage Apple\u2019s Hardened Runtime, ensuring these environment variables are subject to imposed restrictions.[7]Admins can add restrictions to applications by setting the setuid and/or setgid bits, use entitlements, or have a __RESTRICT segment in the Mach-O binary.",
    "Enterprise | T1562 | .003 | Impair Defenses:Impair Command History Logging | Make sure that theHISTCONTROLenvironment variable is set to \"ignoredups\" instead of \"ignoreboth\" or \"ignorespace\".",
    "Enterprise | T1490 | Inhibit System Recovery | Consider technical controls to prevent the disabling of services or deletion of files involved in system recovery. Additionally, ensure that WinRE is enabled using the following command:reagentc /enable.[8]",
    "Enterprise | T1036 | .007 | Masquerading:Double File Extension | Disable the default to \"hide file extensions for known file types\" in Windows OS.[9][10]",
    "Enterprise | T1556 | Modify Authentication Process | Ensure only valid password filters are registered. Filter DLLs must be present in Windows installation directory (C:\\Windows\\System32\\by default) of a domain controller and/or local computer with a corresponding entry inHKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Control\\Lsa\\Notification Packages.Starting in Windows 11 22H2, theEnableMPRNotificationspolicy can be disabled through Group Policy or through a configuration service provider to prevent Winlogon from sending credentials to network providers.[11]",
    " |  | .002 | Password Filter DLL | Ensure only valid password filters are registered. Filter DLLs must be present in Windows installation directory (C:\\Windows\\System32\\ by default) of a domain controller and/or local computer with a corresponding entry in HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Control\\Lsa\\Notification Packages.",
    " |  | .008 | Network Provider DLL | Starting in Windows 11 22H2, theEnableMPRNotificationspolicy can be disabled through Group Policy or through a configuration service provider to prevent Winlogon from sending credentials to network providers.[11]",
    "Enterprise | T1135 | Network Share Discovery | Enable Windows Group Policy \"Do Not Allow Anonymous Enumeration of SAM Accounts and Shares\" security setting to limit users who can enumerate network shares.[12]",
    "Enterprise | T1003 | OS Credential Dumping | Consider disabling or restricting NTLM.[13]Consider disabling WDigest authentication.[14]",
    " |  | .001 | LSASS Memory | Consider disabling or restricting NTLM.[13]Consider disabling WDigest authentication.[14]",
    " |  | .002 | Security Account Manager | Consider disabling or restricting NTLM.[13]",
    " |  | .005 | Cached Domain Credentials | Consider limiting the number of cached credentials (HKLM\\SOFTWARE\\Microsoft\\Windows NT\\Current Version\\Winlogon\\cachedlogonscountvalue)[15]",
    "Enterprise | T1542 | .005 | Pre-OS Boot:TFTP Boot | Follow vendor device hardening best practices to disable unnecessary and unused features and services, avoid using default configurations and passwords, and introduce logging and auditing for detection.",
    "Enterprise | T1563 | .002 | Remote Service Session Hijacking:RDP Hijacking | Change GPOs to define shorter timeouts sessions and maximum amount of time any single session can be active. Change GPOs to specify the maximum amount of time that a disconnected session stays active on the RD session host server.[16]",
    "Enterprise | T1021 | .001 | Remote Services:Remote Desktop Protocol | Change GPOs to define shorter timeouts sessions and maximum amount of time any single session can be active. Change GPOs to specify the maximum amount of time that a disconnected session stays active on the RD session host server.[16]",
    "Enterprise | T1053 | Scheduled Task/Job | Configure settings for scheduled tasks to force tasks to run under the context of the authenticated account instead of allowing them to run as SYSTEM. The associated Registry key is located atHKLM\\SYSTEM\\CurrentControlSet\\Control\\Lsa\\SubmitControl. The setting can be configured through GPO: Computer Configuration > [Policies] > Windows Settings > Security Settings > Local Policies > Security Options: Domain Controller: Allow server operators to schedule tasks, set to disabled.[17]",
    " |  | .002 | At | Configure settings for scheduled tasks to force tasks to run under the context of the authenticated account instead of allowing them to run as SYSTEM. The associated Registry key is located atHKLM\\SYSTEM\\CurrentControlSet\\Control\\Lsa\\SubmitControl. The setting can be configured through GPO: Computer Configuration > [Policies] > Windows Settings > Security Settings > Local Policies > Security Options: Domain Controller: Allow server operators to schedule tasks, set to disabled.[17]",
    " |  | .005 | Scheduled Task | Configure settings for scheduled tasks to force tasks to run under the context of the authenticated account instead of allowing them to run as SYSTEM. The associated Registry key is located at HKLM\\SYSTEM\\CurrentControlSet\\Control\\Lsa\\SubmitControl. The setting can be configured through GPO: Computer Configuration > [Policies] > Windows Settings > Security Settings > Local Policies > Security Options: Domain Controller: Allow server operators to schedule tasks, set to disabled.[17]",
    "Enterprise | T1553 | Subvert Trust Controls | Windows Group Policy can be used to manage root certificates and theFlagsvalue ofHKLM\\SOFTWARE\\Policies\\Microsoft\\SystemCertificates\\Root\\ProtectedRootscan be set to 1 to prevent non-administrator users from making further root installations into their own HKCU certificate store.[18]",
    " |  | .004 | Install Root Certificate | Windows Group Policy can be used to manage root certificates and theFlagsvalue ofHKLM\\SOFTWARE\\Policies\\Microsoft\\SystemCertificates\\Root\\ProtectedRootscan be set to 1 to prevent non-administrator users from making further root installations into their own HKCU certificate store.[18]",
    "Enterprise | T1552 | Unsecured Credentials | There are multiple methods of preventing a user's command history from being flushed to their .bash_history file, including use of the following commands:set +o historyandset -o historyto start logging again;unset HISTFILEbeing added to a user's .bash_rc file; andln -s /dev/null ~/.bash_historyto write commands to/dev/nullinstead.",
    " |  | .003 | Shell History | There are multiple methods of preventing a user's command history from being flushed to their .bash_history file, including use of the following commands:set +o historyandset -o historyto start logging again;unset HISTFILEbeing added to a user's .bash_rc file; andln -s /dev/null ~/.bash_historyto write commands to/dev/nullinstead.In Zsh,fc -pcan be used to create a private history session. However, previous history will be unavailable to the user until the session ends. Usingunset HISTFILEand writing commands to/dev/nullcan also be used, similarly to Bash.In PowerShell, users can utilizeSet-PSReadLineOptionto modify how commands are saved into history. Setting-HistorySaveStyle SaveNothingprevents command history from being saved onto the file. Note that setting it fromSaveNothingtoSaveIncrementallyin the same session will cause all commands from that session to be saved. Alternatively,-AddToHistoryHandlercan be used to filter certain commands from being saved into the history file."
]